# systemd/castctld.service
[Unit]
Description=Chromecast Control

[Service]
EnvironmentFile=/etc/environment
ExecStart=/home/gavin/bin/castctld
Restart=always

[Install]
WantedBy=default.target
